# email-notification.timer
[Unit]
Description=check for any new emails

[Timer]
OnBootSec=0
OnCalendar=*-*-* 08..20:0/2:00
